# deploy/ilvcuisine.service
[Unit]
Description=ILV Cuisine (Flask/Gunicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ilvcuisine
Environment="PATH=/home/ubuntu/ilvcuisine/venv/bin"
ExecStart=/home/ubuntu/ilvcuisine/venv/bin/gunicorn \
    --workers 2 --threads 4 --timeout 120 \
    --bind 127.0.0.1:3024 --bind 172.17.0.1:3024 \
    app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
